# fast_api_experiment/fast_api_experiment.service
[Unit]
Description=Hypercorn instance for Fast API experiment
After=network.target

[Service]
User=nginx
Group=nginx
WorkingDirectory=/var/www/modules/fast_api_experiment
Environment="PATH=/var/www/modules/fast_api_experiment"
ExecStart=/usr/local/bin/hypercorn fast_api_experiment.py --reload --config fast_api_experiment.toml

[Install]
WantedBy=multi-user.target